# mlflow.service
[Unit]
Description=Start MLflow

[Service]
Type=oneshot
ExecStart=/mlflow/bin/launch.sh

[Install]
WantedBy=multi-user.target